# PiServer/http_server.service
[Unit]
Description=Launch HTTP server
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/usr/bin/python /home/pi/server.py
Restart=always
Type=simple
RestartSec=10
KillMode=process
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target